# bin/systemd/gunicorn.gaborone.socket
[Unit]
Description=gunicorn socket (gaborone)

[Socket]
ListenStream=/run/gunicorn.gaborone.sock

[Install]
WantedBy=sockets.target