# bin/systemd/gunicorn.capetown.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.capetown.sock

[Install]
WantedBy=sockets.target