# bin/systemd/gunicorn-kampala.socket
[Unit]
Description=gunicorn socket (kampala)

[Socket]
ListenStream=/run/gunicorn-kampala.sock

[Install]
WantedBy=sockets.target