# bin/systemd/gunicorn.capetown.socket
[Unit]
Description=gunicorn socket (capetown)

[Socket]
ListenStream=/run/gunicorn.capetown.sock

[Install]
WantedBy=sockets.target