# bin/systemd/gunicorn-blantyre.service
[Unit]
Description=gunicorn daemon
Requires=gunicorn-blantyre.sock
After=network.target

[Service]
User=ambition
Group=www-data
WorkingDirectory=/home/ambition/ambition
ExecStart=/home/ambition/.venvs/ambition/bin/gunicorn \
    --access-logfile - \
    --workers 3 \
    --bind unix:/run/gunicorn-blantyre.sock \
    ambition.wsgi.blantyre:application
SyslogIdentifier=gunicorn-blantyre
Restart=always


[Install]
WantedBy=multi-user.target