# bin/systemd/gunicorn-kampala.service
[Unit]
Description=gunicorn daemon
Requires=gunicorn-kampala.sock
After=network.target

[Service]
User=ambition
Group=www-data
WorkingDirectory=/home/ambition/ambition
ExecStart=/home/ambition/.venvs/ambition/bin/gunicorn \
    --access-logfile - \
    --workers 3 \
    --bind unix:/run/gunicorn-kampala.sock \
    ambition.wsgi.kampala:application
SyslogIdentifier=gunicorn-kampala
Restart=always

[Install]
WantedBy=multi-user.target