# bin/systemd/gunicorn-harare.socket
[Unit]
Description=gunicorn socket (harare)

[Socket]
ListenStream=/run/gunicorn-harare.socket

[Install]
WantedBy=sockets.target