# bin/systemd/gunicorn-capetown.socket
[Unit]
Description=gunicorn socket (capetown)

[Socket]
ListenStream=/run/gunicorn-capetown.socket

[Install]
WantedBy=sockets.target