# bin/systemd/gunicorn.blantyre.socket
[Unit]
Description=gunicorn socket (blantyre)

[Socket]
ListenStream=/run/gunicorn-blantyre.socket

[Install]
WantedBy=sockets.target