# bin/systemd/gunicorn.lilongwe.socket
[Unit]
Description=gunicorn socket (lilongwe)

[Socket]
ListenStream=/run/gunicorn-lilongwe.socket

[Install]
WantedBy=sockets.target