# bin/systemd/gunicorn.harare.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.harare.sock

[Install]
WantedBy=sockets.target