# bin/systemd/gunicorn.socket
[Unit]
Description=gunicorn socket (blantyre)

[Socket]
ListenStream=/run/gunicorn-blantyre.sock

[Install]
WantedBy=sockets.target