# bin/systemd/gunicorn.lilongwe.socket
[Unit]
Description=gunicorn socket (lilongwe)

[Socket]
ListenStream=/run/gunicorn.lilongwe.sock

[Install]
WantedBy=sockets.target