# bin/systemd/gunicorn.kampala.socket
[Unit]
Description=gunicorn socket (kampala)

[Socket]
ListenStream=/run/gunicorn-kampala.socket

[Install]
WantedBy=sockets.target