# bin/systemd/gunicorn.harare.socket
[Unit]
Description=gunicorn socket (harare)

[Socket]
ListenStream=/run/gunicorn-harare.sock

[Install]
WantedBy=sockets.target