# bin/systemd/gunicorn-gaborone.service
[Unit]
Description=gunicorn daemon
Requires=gunicorn-gaborone.sock
After=network.target

[Service]
User=ambition
Group=www-data
WorkingDirectory=/home/ambition/ambition
ExecStart=/home/ambition/.venvs/ambition/bin/gunicorn \
    --access-logfile - \
    --workers 3 \
    --bind unix:/run/gunicorn-gaborone.sock \
    ambition.wsgi.gaborone:application

[Install]
WantedBy=multi-user.target